# build/03-supervisor/dnsmasq.conf
# Start dnsmasq server.
[program:dnsmasq]
command=/usr/local/bin/dns-init
priority=100
