# build/03-supervisor/kick-off.conf
# Initialization scripts.
[program:kick-off]
command=/usr/local/bin/kick-off
priority=1
startsecs=0
